# scripts/tripvar.service
[Unit]
Description=TripVar Application Stack
Requires=docker.service
After=docker.service
StartLimitInterval=0

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/tripvar
ExecStart=/usr/bin/docker-compose -f /opt/tripvar/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker-compose -f /opt/tripvar/docker-compose.prod.yml down
ExecReload=/usr/bin/docker-compose -f /opt/tripvar/docker-compose.prod.yml restart
TimeoutStartSec=0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target